# checking_myIP.service
[Unit]
Description=My Script Service

[Service]
Type=oneshot
ExecStart=/home/kali/scripts/know_ch_ip.py

