# checking_myIP.timer
[Unit]
Description=Run My Script every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
